# resources/supervisor/cosmian_helloworld.ini
[program:cosmian_helloworld]
command=cosmian_helloworld
directory=/usr/sbin
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/cosmian_helloworld.err.log
stdout_logfile=/var/log/cosmian_helloworld.out.log